# Linux-Ubuntu-Script/laravel-install-Scripts/laravelStart.service
[Unit]
Description=Laravel Service
After=network.target

[Service]
ExecStart=/home/admin/laravel-start.sh

[Install]
WantedBy=multi-user.target